# archive/auto-rotate-touchscreen.service
[Unit]
Description=Auto-rotate touchscreen input to match display orientation
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/mars/dotfiles/scripts/auto-rotate-touchscreen.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target